# managment/bt-wifi-config.service
[Unit]
Description=Bluetooth WiFi Configuration Service using SPP
Requires=bluetooth.service systemd-user-sessions.service network-pre.target
After=bluetooth.service systemd-user-sessions.service network-pre.target
# Wants=network-online.target # Optional: wait for network if the script itself needs it initially

[Service]
User=<your_username>
Group=<your_username> # Often the same as username, e.g., 'pi'
WorkingDirectory=/home/<your_username>/lb
ExecStart=/bin/bash /home/<your_username>/lb/start_bt_wifi.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal # Logs stdout to systemd journal
StandardError=journal # Logs stderr to systemd journal

[Install]
WantedBy=multi-user.target
